# supervisor.conf.d/scribe.ini
[program:facebook-scribe-docker]
command=bash -c 'cd /data/app/facebook-scribe-docker && make foreground-start'
numprocs=1
process_name=scribe
killasgroup=true
autorestart=true
user=root
stdout_logfile=/data/log/supervisor/stdout.log
stderr_logfile=/data/log/supervisor/stderr.log